# supervisor/myprog.ini
[program:myprog]
user=andre
command=/usr/local/bin/python /opt/myprog/bin/bootstrap
process_name=%(program_name)s
stopasgroup=true
killasgroup=true
startsecs=10
stderr_logfile=/var/log/myprog/myprog.stderr.log
minfds=65536
